# debian/valkey.prerm
#!/bin/bash
set -e

# Stop valkey services before package removal
systemctl stop valkey.service
systemctl disable valkey.service

systemctl stop valkey-sentinel.service
systemctl disable valkey-sentinel.service
exit 0

